# recipes-core/buttond/files/buttond.service
[Unit]
Description=Button listener service

[Service]
ExecStart=/usr/bin/buttons
ExecStop=killall -9 buttons

[Install]
WantedBy=multi-user.target
