# recipes-support/rsync/files/rsync.service
[Unit]
Description=rsync daemon

[Service]
Type=forking
ExecStart=/usr/bin/rsync --daemon

[Install]
WantedBy=multi-user.target
